# supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D

[program:mongodb]
command=/mongodb/bin/mongod
autostart=true
autorestart=true
startsecs=2
user=root
